# slackware/rc-local.service
[Unit]
Description=rc.local
 
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/etc/rc.local
 
[Install]
WantedBy=multi-user.target
